# model_server/fastchat_server/conf/controller.conf
[program:controller]
command=python3 -m fastchat.serve.controller
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
