# model_server/fastchat_server/conf/worker.conf
[program:worker]
command=python3 -m fastchat.serve.model_worker --model-path lmsys/vicuna-7b-v1.5
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
