# munchkin/support-files/supervisor/celery.conf
[program:celery]
command=celery -A munchkin worker -B --loglevel=info --pool threads
directory=/apps
user=root
autostart=true
autorestart=true
redirect_stderr=true
