# munchkin/support-files/supervisor/munchkin.conf
[program:munchkin]
command=daphne -b 0.0.0.0 -p 8000 munchkin.asgi:application
directory=/apps
user=root
autostart=true
autorestart=true
redirect_stderr=true
