# model_server/classicfy_aiops_server/conf/service.conf
[program:service]
command=python3 ./server.py
directory=/apps
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0