# depend/saltstack_server/config/supervisor/salt-api.conf
[program:salt-api]
command= /usr/bin/salt-api
autostart=true
priority=10
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/salt-api.log