# depend/saltstack_server/config/supervisor/salt-master.conf
[program:salt-master]
command= /usr/bin/salt-master
autostart=true
priority=10
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/salt-master.log