# model_server/saltstack_server/config/supervisor/salt-minion.conf
[program:salt-minion]
command= /usr/bin/salt-minion
autostart=true
priority=20
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/salt-minion.log